# container_root/etc/supervisord.d/fluxbox.conf
[program:fluxbox]
command=fluxbox
autorestart=true
