# container_root/etc/supervisord.d/websockify.conf
[program:websockify]
command=websockify --web /usr/share/novnc %(ENV_NOVNC_PORT)s localhost:5900
autorestart=true
