# container_root/etc/supervisord.d/wechat-monitor.conf
[program:wechat-monitor]
command=wechat-monitor
autorestart=true