# container_root/etc/supervisord.d/x11vnc.conf
[program:x11vnc]
command=x11vnc -forever -shared
autorestart=true
