# container_root/etc/supervisord.d/xvfb.conf
[program:xvfb]
command=Xvfb %(ENV_DISPLAY)s -screen 0 "%(ENV_VNC_DISPLAY_WIDTH)s"x"%(ENV_VNC_DISPLAY_HEIGHT)s"x24 -listen tcp -ac
autorestart=true
